# supervisor_fpm.conf
[program:fpm]
command=/usr/sbin/php-fpm7.0 -F
autorestart=true
user=root

